# etc/cron.d/cloudradar
# smart
#  m    h   dom  month  dow
   */15 *   *   *    * /opt/cloudradar/smart.sh
