# etc/cron.d/cloudradar
# smart
# m h dom month dow user cmd
* * * * * root /opt/cloudradar/smart.sh
